# infra/opspocket-waitlist.service
[Unit]
Description=OpsPocket Cloud waitlist endpoint
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/opspocket/waitlist-server.py
Restart=always
RestartSec=3
User=root
Group=root
ProtectSystem=strict
ReadWritePaths=/var/lib/opspocket
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
