# infra/scripts/opspocket-snapshot.timer
[Unit]
Description=OpsPocket — daily Hetzner snapshot timer
Documentation=https://github.com/findgriff/opspocket/blob/main/infra/README.md

[Timer]
# 04:00 UTC every day. Randomised by 5 min so bursts of Hetzner API calls
# across customers don't pile onto the same second.
OnCalendar=*-*-* 04:00:00 UTC
RandomizedDelaySec=5min
# If the box was off at 04:00, run as soon as it boots.
Persistent=true
Unit=opspocket-snapshot.service

[Install]
WantedBy=timers.target
