# infra/tenant/opspocket-heartbeat.service
[Unit]
Description=OpsPocket heartbeat — push tenant telemetry to the platform
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/opspocket-heartbeat.sh
# Don't spam logs on transient network failures
StandardError=journal
StandardOutput=journal
# Harden
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
